# docker/crontab/basecron
#############################################################################
#
# This is just a basic cron, it will never be more advnaced than this.
#
#############################################################################
* * * * * php artisan schedule:run
